# ansible/roles/semaphore-controller/templates/semaphore.service.j2
[Unit]
Description=Semaphore UI - Ansible Automation Platform
Documentation=https://docs.semaphoreui.com
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User={{ ansible_user_name }}
Group={{ ansible_user_name }}

# Load environment variables
EnvironmentFile=/etc/default/semaphore

# Working directory (infrastructure repository)
WorkingDirectory={{ infra_repo_path }}

# Semaphore binary with config path
ExecStart={{ semaphore_install_path }}/semaphore server --config {{ semaphore_config_path }}/config.json

# Restart on failure
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ infra_repo_path }} /tmp {{ ansible_user_home }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=semaphore

[Install]
WantedBy=multi-user.target
